# deploy/pidsim-api.service
# systemd unit for the FastAPI backend (gunicorn managing uvicorn workers).
#
# gunicorn is used (rather than bare uvicorn) so the admin "Deploy" action can
# gracefully reload via SIGHUP to the master PID: new workers spawn before old
# ones retire, so the HTTP request that triggered the deploy completes normally.
# That reload needs no sudo because the web process and gunicorn run as the same
# user. A hard `systemctl restart` is only needed for code that can't hot-reload
# (rare); the sudoers rule in deploy/README covers that fallback.
#
# Install:
#   sudo cp /opt/robot/pidsim/deploy/pidsim-api.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now pidsim-api

[Unit]
Description=Matamoe Line Follower PID Simulator API
After=network.target

[Service]
Type=simple
User=matamoe
Group=matamoe
WorkingDirectory=/opt/robot/pidsim/backend
Environment=PIDSIM_DATA_DIR=/opt/robot/pidsim/data
ExecStart=/opt/robot/pidsim/backend/.venv/bin/gunicorn app.main:app \
    --workers 2 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind 127.0.0.1:8700 \
    --timeout 120
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
